# conf/healthagent-restart.timer
[Unit]
Description=Weekly restart of healthAgent to prevent OOM memory leak

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
