# conf/juicefs-backup.service
[Unit]
Description=JuiceFS metadata backup to GCS
After=redis-server.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/juicefs-backup
TimeoutStartSec=10min
PrivateTmp=true
