# conf/juicefs-gc-compact.timer
[Unit]
Description=Weekly JuiceFS gc --compact

[Timer]
OnCalendar=Sun 04:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
